# build/mattermost/files/mattermost.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="minidlna">

<service 
	name="oep/social/mattermost"
	type="service" 
	version="1">

        <create_default_instance 
        	enabled="false" />

        <dependency
            name='network'
            grouping='require_all'
            restart_on='none'
            type='service'>
            <service_fmri value='svc:/milestone/network:default' />
        </dependency>

        <dependency
            name='filesystem'
            grouping='require_all'
            restart_on='none'
            type='service'>
            <service_fmri value='svc:/system/filesystem/local:default' />
        </dependency>

        <dependency
            name='config-file'
            grouping='require_all'
            restart_on='refresh'
            type='path'>
            <service_fmri value='file://localhost/etc/opt/oep/mattermost/config.json' />
        </dependency>

        <method_context>
                <method_credential
                        user='mattermost'
                        group='mattermost' />
        </method_context>

        <exec_method 
		type="method" 
		name="start" 
		exec="/opt/oep/mattermost/bin/platform -c /etc/opt/oep/mattermost/config.json &amp;"
		timeout_seconds="60" />

        <exec_method
		type="method"
		name="stop"
		exec=":kill"
		timeout_seconds="60" />

        <stability 
        	value="Unstable" />

        <template>
            <common_name>
                <loctext xml:lang="C">
                Mattermost - All your team communication in one place, instantly searchable and accessible anywhere.
                </loctext>
            </common_name>
        </template>
</service>
</service_bundle>
